# apps/resilio-sync/resilio-sync-helene.container
[Unit]
Description=Resilio Sync - a fast, reliable file and folder syncing solution for homes and businesses.
After=network-online.target
Wants=network-online.target

[Service]
TimeoutStartSec=900
Restart=always
RestartSec=3

[Container]
AutoUpdate=registry
ContainerName=resilio-sync-helene
Image=ghcr.io/auricom/resilio-sync:3.0.3@sha256:791a35477c7d6e93f2ee51f9a9f718548a59339875ca3a241f76c3e4f6daa7d6
Volume=/usr/etc/resilio-sync/helene.conf:/app/sync.conf:ro
Volume=resilio-sync-helene:/config:rslave
Volume=/var/mnt/vol1/shared-documents:/sync/shared-documents:rslave
Volume=/var/mnt/vol1/photo:/sync/photo:rslave
Volume=/var/mnt/vol1/home/helene:/sync/helene:rslave
User=root:root
SecurityLabelType=container_runtime_t
Label=traefik.enable="true"
Label=traefik.http.routers.resilio-sync-helene.rule=Host'(`resilio-sync-helene.feisar.ovh`)'
Label=traefik.http.routers.resilio-sync-helene.entrypoints="websecure"
Label=traefik.http.services.resilio-sync-helene.loadbalancer.server.port="8888"

[Install]
WantedBy=multi-user.target
